# services/hs-sshgateway/run
#!/bin/sh
echo starting
exec 2>&1
exec /usr/lib/hs/sshgateway/ssh-gateway.py
